# crontab
MAILTO="joel" 
05 17 * * * /usr/bin/rsync -avzP --exclude-from '/home/joel/exclude_me.txt' --delete $HOME/.??* backup@slacktop:/home/backup/backups

